# bigbluebutton-server/systemd/bbb-rap-resque-worker.service.d/lowprio.conf
# Drop-in: lower the priority of BBB's recording-pipeline worker so it can
# never starve live-meeting services (FreeSWITCH, bbb-webrtc-sfu) of CPU.
# Under contention cgroup v2 gives this worker only ~1/100th the CPU share
# of services running with default CPUWeight=100; when the box is idle it
# still gets 100% so recording processing finishes within the 5.5h window.
[Service]
CPUWeight=10
IOWeight=10
Nice=19
IOSchedulingClass=idle
