# bigbluebutton-server/systemd/bbb-webrtc-sfu.service.d/highprio.conf
# Drop-in: give the WebRTC SFU (video forwarding for live meetings) a CPU
# weight 10x the default. Same reasoning as FreeSWITCH — video freezes /
# blurring during class is most often this process being CPU-starved by
# the recording pipeline.
[Service]
CPUWeight=1000
IOWeight=1000
