# bigbluebutton-server/systemd/freeswitch.service.d/highprio.conf
# Drop-in: give FreeSWITCH (audio mixing for live meetings) a CPU weight
# 10x the default. With cgroup v2 this means under contention it gets ~10x
# the CPU share of any default-weight service and ~100x the rap-worker.
# Audio is the most timing-sensitive part of a live meeting; jitter here
# is what users perceive as "glitchy voice".
[Service]
CPUWeight=1000
IOWeight=1000
